# Jamulus@.service.d/Jamulus-Directory4.conf
[Unit]
Description=Directory on 22424

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
