# Jamulus@.service.d/Jamulus-server.conf
[Unit]
Description=Jamulus at drealm.info
Conflicts=Jamulus@original-music.service
Before=Jamulus@original-music.service

[Service]
#CPUSchedulingPriority=90
TimeoutStopSec=2
Nice=-10
