# Jamulus@.service.d/Jamulus-wip.conf
[Unit]
Description=drealm.info Test Server

[Install]
WantedBy=

[Service]
User=peter
Group=peter
CPUSchedulingPriority=99
