# Jamulus@.service.d/Jamulus-private.conf
[Unit]
Description=Private Server at drealm.info
Conflicts=Jamulus@wip.service
Before=Jamulus@wip.service

[Service]
#CPUSchedulingPriority=90
TimeoutStopSec=2
Nice=-10
