# Jamulus@.service.d/Jamulus-Directory7.conf
[Unit]
Description=Directory on 22724

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
