# Jamulus@.service.d/Jamulus-Central3.conf
[Unit]
Description=Central Server List (22324)

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
