# Jamulus@.service.d/Jamulus-Central1.conf
[Unit]
Description=Directory Server (22124)

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
