# Jamulus@.service.d/Jamulus-wip.conf
[Unit]
Description=drealm.info Test Server
Conflicts=Jamulus@private.service
After=Jamulus@private.service

[Install]
WantedBy=

[Service]
User=peter
Group=peter
