# Jamulus@.service.d/Jamulus-server.conf
[Unit]
Description=Jamulus at drealm.info

[Service]
CPUSchedulingPriority=99
TimeoutStopSec=2
