# Jamulus@.service.d/Jamulus-Directory6.conf
[Unit]
Description=Directory on 22624

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
