# Jamulus@.service.d/Jamulus-Central6.conf
[Unit]
Description=Central Server List (22624)

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
