# Jamulus@.service.d/Jamulus-Directory3.conf
[Unit]
Description=Directory on 22324

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
