# Jamulus@.service.d/Jamulus-Central1.conf
[Unit]
Description=Central Server List (22524)

[Service]
CPUSchedulingPolicy=other
CPUSchedulingPriority=0
Nice=-5
