# Jamulus@.service.d/Jamulus-private.conf
[Unit]
Description=Private Server at drealm.info

[Service]
#CPUSchedulingPriority=90
TimeoutStopSec=2
Nice=-10
